# systemd/cam-watchdog.timer
[Unit]
Description=Kamera kontrol zamanlayıcı

[Timer]
OnCalendar=*-*-* *:*:00/1m
AccuracySec=30s

[Install]
WantedBy=timers.target