# systemd/cam-watchdog.timer
[Unit]
Description=Kamera kontrol zamanlayıcı

[Timer]
OnCalendar=*-*-* *:*:00/5m
AccuracySec=1m

[Install]
WantedBy=timers.target