# systemd/egg-counter.service
[Unit]
Description=Yumurta Sayıcı Ana Uygulama
After=network.target
Wants=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/kaplan/Desktop/Azim-Tav/Yumurta_sayıcı
ExecStart=/usr/bin/python3 /home/kaplan/Desktop/Azim-Tav/Yumurta_sayıcı/main.py --headless
Restart=on-failure
RestartSec=5
Environment=PYTHONUNBUFFERED=1
Environment=OPENCV_FFMPEG_CAPTURE_OPTIONS=threads;1|fflags;nobuffer
Environment=OPENCV_FFMPEG_MULTITHREADED=0
# Hata durumunda otomatik tamir: watchdog
WatchdogSec=30
# Kamera hatası gibi durumlarda yeniden başlat
RestartPreventExitStatus=1 2 3 4 5 6 7 8 9 10
# Loglama
StandardOutput=journal
StandardError=journal
# CPU/Memory limitleri (RPi5 için optimize)
CPUQuota=80%
MemoryLimit=1G

[Install]
WantedBy=multi-user.target