# systemd/health-check.service
[Unit]
Description=Health Check and Auto-Repair Service
After=network.target
Wants=network.target

[Service]
Type=oneshot
User=pi
WorkingDirectory=/home/kaplan/Desktop/Azim-Tav/Yumurta_sayıcı
ExecStart=/bin/bash /home/kaplan/Desktop/Azim-Tav/Yumurta_sayıcı/health_check.sh
RemainAfterExit=no
# Her 5 dakikada bir çalıştır
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target