# systemd/update.service
[Unit]
Description=Update and Restart Services
After=network.target
Wants=network.target

[Service]
Type=oneshot
User=pi
WorkingDirectory=/home/kaplan/Desktop/Azim-Tav/Yumurta_sayıcı
ExecStart=/bin/bash /home/kaplan/Desktop/Azim-Tav/Yumurta_sayıcı/update_and_restart.sh
RemainAfterExit=no
# Hata durumunda logla ama devam et
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target